# init/ddns.service
[Unit]
Description=DDNS
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=r53

[Install]
WantedBy=timers.target
